# ansible/pulse/templates/pulse.service.j2
[Unit]
Description=OCTALUM-PULSE - Your Infrastructure's Heartbeat
Documentation=https://docs.pulse.harery.com
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart={{ pulse_install_dir }}/pulse --auto
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
