# deployments/systemd/pulse.timer
[Unit]
Description=Weekly OCTALUM-PULSE Maintenance
Documentation=https://docs.pulse.harery.com

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
